# internal/systemd/templates/update.service.tmpl
[Unit]
Description=DayZ scheduled update (dayzctl)
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/dayzctl update
TimeoutStartSec=3600
